# templates/couchpotato.service.j2
{{ ansible_managed | comment }}
[Service]
SupplementaryGroups={{ couchpotato_systemd_supplementary_groups | join(' ') }}
